# borg-warn-if-no-recent-backup.timer
[Unit]
Description=Timer for borg-warn-if-no-recent-backup

[Timer]
OnCalendar=daily
RandomizedDelaySec=2h

[Install]
WantedBy=timers.target
