# borg-warn-if-no-recent-backup.timer
[Unit]
Description=Timer for borg-warn-if-no-recent-backup

[Timer]
OnCalendar=0/6:00:00
RandomizedDelaySec=2h

[Install]
WantedBy=timers.target
