# borg-backup.service
[Unit]
Description=Do borg backup

[Service]
User=root
StandardError=file:/var/log/borg.log
ExecStart=/bin/bash /bin/borg-backup
