# borg-warn-if-no-recent-backup.service
[Unit]
Description=Warn if no recent borg backups are found

[Service]
User=root
ExecStart=/bin/borg-warn-if-no-recent-backup
